# scripts/startup/openrc/falcon-as
#!/sbin/openrc-run
# Falcon HTTP Application Server - OpenRC Init Script

name="falcon-as"
description="Falcon HTTP Application Server"

command="/usr/local/bin/falcon-as"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
    after firewall
}

start_pre() {
    # set ulimit for open files
    ulimit -n 1000000

    # set transparent hugepages
    sysctl -w vm.nr_hugepages=1024

    # Ensure required directories exist
    checkpath --directory --owner falcon-http:falcon-http --mode 0750 /etc/falcon-http
    checkpath --directory --owner falcon-http:falcon-http --mode 0775 /var/www
}
